# usr/lib/systemd/system/spice-vdagentd.service.d/10-cernvm.conf
[Unit]
ConditionPathExists=/dev/virtio-ports/com.redhat.spice.0
