# usr/lib/systemd/system/cloud-final.service.d/10-cernvm.conf
[Unit]
ConditionPathExists=/mnt/.rw/context/ucontext
Before=multi-user.target
