# usr/lib/systemd/system/getty@tty1.service.d/decontaminate-loginprompt.conf
[Service]
ExecStartPre=-/bin/kill -55 1
