# usr/lib/systemd/system/cernvm-contextualized.target
[Unit]
Description=CernVM contextualization fetched
Wants=cernvm-online.service
