# usr/lib/systemd/system/xrdp.service.d/10-cernvm.conf
[Service]
ExecStartPre=/etc/cernvm/config -X
